# configuration/supervisor-app.conf

[supervisord]
nodaemon=true

[program:uwsgi]
command = /usr/sbin/uwsgi --ini /home/docker/code/uwsgi.ini

[program:nginx]
command = /usr/sbin/nginx